# lib/systemd/system/LCD-Stats.service
[Unit]
Description=LCD-Stats

[Service]
ExecStart=/bin/bash /usr/local/LCD-Stats/LCD-Stats.sh

[Install]
WantedBy=multi-user.target
